# systemd/librespot.service
[Unit]
Description=Spotify Connect Receiver
After=network-online.target sound.target
Wants=network-online.target

[Service]
Type=simple
User=mopidy
Group=audio
EnvironmentFile=-/etc/default/spotify-panel
Environment=SPOTIFY_PANEL_PYTHONPATH=/opt/spotify-panel/src
ExecStartPre=+/usr/bin/install -d -o mopidy -g audio -m 775 /opt/spotify-panel
ExecStartPre=+/usr/bin/install -d -o mopidy -g audio -m 775 /opt/spotify-panel/covers
ExecStartPre=+/usr/bin/install -d -o mopidy -g audio -m 775 /var/cache/librespot
ExecStart=/usr/local/bin/librespot \
  --name "${LIBRESPOT_NAME}" \
  --backend alsa \
  --device "${LIBRESPOT_DEVICE}" \
  --bitrate "${LIBRESPOT_BITRATE}" \
  --initial-volume "${LIBRESPOT_INITIAL_VOLUME}" \
  --cache "${LIBRESPOT_CACHE}" \
  --disable-audio-cache \
  --onevent /usr/local/bin/spotify-event.sh
Restart=on-failure
RestartSec=5
TimeoutStopSec=15
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/opt/spotify-panel /var/cache/librespot

[Install]
WantedBy=multi-user.target
