# systemd/librespot.service
[Unit]
Description=Spotify Connect Receiver
After=network-online.target sound.target
Wants=network-online.target

[Service]
User=mopidy
Group=audio
EnvironmentFile=-/etc/default/spotify-panel
# Make sure the event hook can write runtime files even after manual cleanup.
ExecStartPre=+/usr/bin/install -d -o mopidy -g audio -m 775 /opt/spotify-panel
ExecStartPre=+/usr/bin/install -d -o mopidy -g audio -m 775 /var/cache/librespot
ExecStart=/usr/local/bin/librespot \
  --name "${LIBRESPOT_NAME}" \
  --backend alsa \
  --device "${LIBRESPOT_DEVICE}" \
  --bitrate "${LIBRESPOT_BITRATE}" \
  --initial-volume "${LIBRESPOT_INITIAL_VOLUME}" \
  --cache "${LIBRESPOT_CACHE}" \
  --disable-audio-cache \
  --onevent /usr/local/bin/spotify-event.sh
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
